# physical/physical/files/lvshost/etc/network/if-up.d/iptables
#!/bin/bash

if ! iptables -t nat -C POSTROUTING -o eth0 -j MASQUERADE ; then
        iptables -t nat -A POSTROUTING -o eth0 -j MASQUERADE
fi


