# moriya_develop/ansible/template/web-setup.service.j2
[Unit]
Description=Web network setup
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/work/web_setup.sh

[Install]
WantedBy=multi-user.target

